# .platform/systemd/laravel-schedule.timer
[Unit]
Description=Run laravel scheduled tasks

[Timer]
OnBootSec=1min
OnUnitActiveSec=1m

[Install]
WantedBy=timers.target
